# example.cron
* * * * * python "/home/tpk/project/cron/cex.py"
