# example.cron
* * * * * python "/home/tpk/projects/cron/cex.py"
