# example.cron
# uses local python installation and logs error to c.log
* * * * * ~/anaconda3/bin/python ~/projects/cron/cex.py >> ~/projects/cron/c.log 2>&1
